# etc/profile.d/amdgpu.sh
export AMD_DEBUG=nodcc
export OCL_ICD_VENDORS=mesa
export RADV_PERFTEST=aco
